# scripts/bathyimager.service
[Unit]
Description=BathyImager Seabed Imager Service
Documentation=https://github.com/Mike-Bollinger/BathyImager-Seabed-Imager
After=network.target multi-user.target
Wants=network.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=bathyimager
Group=bathyimager
WorkingDirectory=/home/bathyimager/BathyCat-Seabed-Imager/src
Environment=PYTHONPATH=/home/bathyimager/BathyCat-Seabed-Imager/src
Environment=PYTHONUNBUFFERED=1
ExecStartPre=/bin/sleep 10
ExecStartPre=/bin/bash -c 'chmod 666 /dev/video* 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'chmod 666 /dev/ttyUSB* 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'chmod 666 /dev/ttyACM* 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'chmod 666 /dev/gpiomem 2>/dev/null || true'
ExecStartPre=/bin/bash -c 'chmod 666 /dev/mem 2>/dev/null || true'
ExecStart=/home/bathyimager/BathyCat-Seabed-Imager/run_bathyimager.sh
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=bathyimager

# Security settings (relaxed for hardware access)
NoNewPrivileges=false
ProtectHome=false
ProtectSystem=false
ReadWritePaths=/home/bathyimager/BathyCat-Seabed-Imager /var/log/bathyimager /media /mnt /dev /sys
PrivateTmp=false
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectControlGroups=false
RestrictRealtime=false
RestrictSUIDSGID=false
LockPersonality=false
RemoveIPC=false

# Resource limits
MemoryHigh=512M
MemoryMax=1G
TasksMax=100

# Capability restrictions (allow hardware access)
CapabilityBoundingSet=CAP_DAC_OVERRIDE CAP_FOWNER CAP_CHOWN CAP_SETGID CAP_SETUID CAP_SYS_RAWIO CAP_SYS_ADMIN
AmbientCapabilities=CAP_DAC_OVERRIDE CAP_FOWNER CAP_SYS_RAWIO

# Device access (for camera, GPS, and GPIO) 
DeviceAllow=/dev/video* rw
DeviceAllow=/dev/ttyUSB* rw
DeviceAllow=/dev/ttyACM* rw
DeviceAllow=/dev/serial* rw
DeviceAllow=/dev/gpiomem rw
DeviceAllow=/dev/mem rw
DeviceAllow=char-* rw
DeviceAllow=block-* rw
SupplementaryGroups=gpio video dialout i2c spi audio input
DevicePolicy=auto

# Network restrictions
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# System call filtering (minimal restrictions for hardware access)
# SystemCallFilter=~@obsolete @keyring @module @mount
# SystemCallErrorNumber=EPERM

# File system restrictions (disabled for hardware access)
# ProtectProc=invisible
# ProcSubset=pid

[Install]
WantedBy=multi-user.target
Alias=bathyimager-service.service