# scripts/gps-periodic-sync.timer
[Unit]
Description=BathyCat GPS Periodic Time Synchronization Timer
Documentation=https://github.com/your-org/bathycat-seabed-imager
Requires=gps-periodic-sync.service

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target