# scripts/gps-periodic-sync.service
[Unit]
Description=BathyCat GPS Periodic Time Synchronization
Documentation=https://github.com/your-org/bathycat-seabed-imager
ConditionPathExists=/home/bathyimager/BathyCat-Seabed-Imager/src/gps_time_sync.py

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/home/bathyimager/BathyCat-Seabed-Imager/venv/bin/python /home/bathyimager/BathyCat-Seabed-Imager/src/gps_time_sync.py --periodic
WorkingDirectory=/home/bathyimager/BathyCat-Seabed-Imager
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120

# Environment
Environment=PYTHONPATH=/home/bathyimager/BathyCat-Seabed-Imager/src
Environment=PATH=/home/bathyimager/BathyCat-Seabed-Imager/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Security settings
NoNewPrivileges=false
PrivateTmp=true
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectControlGroups=false
RestrictRealtime=false
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
MemoryLimit=128M
CPUQuota=25%